# docs/inbox.service
[Unit]
Description=Inbox API server
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/misc/inbox
ExecStart=%h/.npm-global/bin/bun inbox serve
Environment=INBOX_PORT=3333
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
